# 99-motu-m4-audio-optimizer.rules
# Debug version for MOTU M4
SUBSYSTEM=="sound", ACTION=="add", \
    RUN+="/bin/logger 'Sound-Add-Event: %k'"

SUBSYSTEM=="sound", KERNEL=="controlC*", ACTION=="add", \
    RUN+="/bin/logger 'ControlC-Event: %k'", \
    RUN+="/bin/sleep 0.5", \
    RUN+="/bin/sh -c 'if ls /proc/asound/*/id 2>/dev/null | xargs cat | grep -q M4; then logger \"MOTU M4 detected - Starting service (delayed)\"; systemctl --no-block start motu-m4-dynamic-optimizer-delayed.service; fi'"

SUBSYSTEM=="sound", KERNEL=="card*", ACTION=="remove", \
    RUN+="/bin/logger 'Sound-Remove: Stopping service'", \
    RUN+="/bin/systemctl --no-block stop motu-m4-dynamic-optimizer.service"
