# motu-m4-dynamic-optimizer-delayed.service
[Unit]
Description=MOTU M4 Dynamic Audio Optimizer (Delayed for User Sessions)
Documentation=man:systemd.service(5)
After=sound.target multi-user.target graphical.target
Wants=sound.target
# Delayed service - wartet auf User-Session-Audio-Services

[Service]
Type=oneshot
# Verzögerte Ausführung für User-Session-Audio-Services
ExecStartPre=/bin/sleep 15
ExecStart=/usr/local/bin/motu-m4-dynamic-optimizer.sh once-delayed
ExecStop=/usr/local/bin/motu-m4-dynamic-optimizer.sh stop
RemainAfterExit=yes
User=root
StandardOutput=journal
StandardError=journal

# Umgebungsvariablen
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Timeout erhöhen für Wartezeiten
TimeoutStartSec=60

# Service läuft nur on-demand, kein automatischer Restart
# Restart=no (default)

[Install]
# Nicht automatisch starten - wird von udev gestartet
# WantedBy=multi-user.target
